# rootfs/etc/s6-overlay/s6-rc.d/ii/run
#!/command/with-contenv sh

echo "-- Connexion au serveur ${SERVER_IP}"
exec s6-setuidgid alpine ii -i /app/servers -s ${SERVER_IP} -p ${SERVER_PORT} -n ${NICKNAME} -f ${NICKNAME}